# bin/plugins/lox-audioserver/lox-audioserver.service
[Unit]
Description=Lox-Audioserver
After=network.target

[Service]
Type=simple
User=loxberry
Group=loxberry
WorkingDirectory=/opt/loxberry/bin/plugins/lox-audioserver/plugins/lox-audioserver/app
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
